# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn blog:web -b localhost:5000
directory=/www
